# nodes/作业/cluster00/redis/redis_6371.ini
[program:predixy]
command=/opt/redis-server/bin/predixy /opt/redis-server/conf/predixy.conf ; port 6370
directory=/opt/redis-server/conf/
autostart=true
autorestart=true
user=www
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/redis-server/logs/predixy.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10

[program:redis_6371]
command=/opt/redis-server/bin/redis-server /opt/redis-server/etc/redis_6371.conf
directory=/opt/redis-server/data/
autostart=true
autorestart=true
user=www
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/redis-server/logs/redis_6371.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
